# ops/systemd/user/classmanager-backup-verify.service
[Unit]
Description=Run classmanager-multi weekly backup verification
Wants=network-online.target
After=network-online.target
OnFailure=classmanager-backup-verify-alert.service

[Service]
Type=oneshot
WorkingDirectory=__PROJECT_DIR__
EnvironmentFile=-__PROJECT_DIR__/.env.runtime
ExecStart=/usr/bin/npm run backup:verify
