# ops/systemd/user/classmanager-backup-freshness.service
[Unit]
Description=Check classmanager backup freshness

[Service]
Type=oneshot
WorkingDirectory=__PROJECT_DIR__
ExecStart=/usr/bin/npm run backup:check-freshness
