# ops/systemd/user/classmanager-backup.service
[Unit]
Description=Run classmanager-multi daily SQLite backup
Wants=network-online.target
After=network-online.target
OnFailure=classmanager-backup-alert.service

[Service]
Type=oneshot
WorkingDirectory=__PROJECT_DIR__
EnvironmentFile=-__PROJECT_DIR__/.env.runtime
Environment=PATH=/opt/node-current/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/env npm run backup:daily
