# ops/systemd/user/classmanager-backup-verify.timer
[Unit]
Description=Weekly timer for classmanager-multi backup verification

[Timer]
OnCalendar=Mon *-*-* 02:30:00
Persistent=true
AccuracySec=1m

[Install]
WantedBy=timers.target
