# ops/systemd/user/classmanager-backup.service
[Unit]
Description=Run classmanager-multi daily SQLite backup
Wants=network-online.target
After=network-online.target
OnFailure=classmanager-backup-alert.service

[Service]
Type=oneshot
WorkingDirectory=__PROJECT_DIR__
ExecStart=/usr/bin/npm run backup:daily
