# ops/systemd/user/classmanager-backup-alert.service
[Unit]
Description=Alert when classmanager backup fails

[Service]
Type=oneshot
WorkingDirectory=/path/to/classmanager-multi
ExecStart=/usr/bin/node scripts/backup-failure-notify.js --unit classmanager-backup.service --reason backup-failed
