# ops/systemd/user/classmanager-github-deploy.timer
[Unit]
Description=Check GitHub for ClassManager updates every 5 minutes

[Timer]
OnBootSec=2min
OnUnitActiveSec=5min
RandomizedDelaySec=15s
Persistent=true
Unit=classmanager-github-deploy.service

[Install]
WantedBy=timers.target
