# ops/systemd/user/classmanager-app.service
[Unit]
Description=ClassManager Multi application
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
WorkingDirectory=__PROJECT_DIR__
EnvironmentFile=-__PROJECT_DIR__/.env.runtime
Environment=NODE_ENV=production
Environment=PORT=3002
Environment=PATH=/opt/node-current/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/usr/bin/env node server.js
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
